# schedule-display.service
[Unit]
Description=Pi Schedule Display
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/schedule-display
ExecStart=/usr/bin/python3 /home/pi/schedule-display/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
